# deployment/systemd/tennis-bot.service
[Unit]
Description=Tennis Discovery Agent - Discord Bot
After=network.target
Documentation=https://github.com/ishidafuu/isdf_tennis_discovery_agent

[Service]
Type=simple
# NOTE: 以下のプレースホルダーはセットアップスクリプトで自動置換されます
# 手動で設置する場合は、実際の値に置き換えてください
#   USER_NAME: 実際のユーザー名 (例: ishidafuu)
#   HOME_DIR: ホームディレクトリ (例: /home/ishidafuu)
#   DOTENVX_PATH: dotenvxの実際のパス (例: /usr/bin/dotenvx または /usr/local/bin/dotenvx)
User=USER_NAME
WorkingDirectory=HOME_DIR/isdf_tennis_discovery_agent
ExecStart=DOTENVX_PATH run -- HOME_DIR/isdf_tennis_discovery_agent/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# セキュリティ設定
NoNewPrivileges=true
PrivateTmp=true

# 環境変数（必要に応じて追加）
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
